# deploy/supervisord.conf
; supervisord config for the multi-model-rag all-in-one image.
; Both processes log to stdout/stderr so Docker / Kubernetes capture them.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=/usr/local/bin/uvicorn multi_model_rag.dashboard_api:app
    --host %(ENV_HOST)s --port %(ENV_BACKEND_PORT)s --proxy-headers
directory=/app
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:dashboard]
command=/usr/bin/node server.js
directory=/app/dashboard
environment=
    HOSTNAME="0.0.0.0",
    PORT="%(ENV_DASHBOARD_PORT)s",
    NEXT_PUBLIC_MULTIMODEL_API_BASE_URL="%(ENV_NEXT_PUBLIC_MULTIMODEL_API_BASE_URL)s",
    NODE_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
